# hioki-app.service
[Unit]
Description=Hioki Data Logger
After=network.target
# Watchdog will restart this service if it hangs and gets killed
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Hioki_data_logger
ExecStart=/home/pi/Hioki_data_logger/.venv/bin/python3 main.py
Restart=always
RestartSec=3

# Allow GUI (needed for Qt on Raspberry Pi desktop)
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

[Install]
WantedBy=graphical.target
